# templates/update-letsencrypt.timer.j2
[Unit]
Description=Trigger Letsencrypt updates

[Timer]
OnCalendar=daily
Persistent=true
